# crontab
*/1	*	*	*	*	/update-dns.sh >> /var/log/update-dns.log
