# crontab
*/5 * * * * root /update-dns.sh >> /var/log/cron.log 2>&1

*/5	*	*	*	*	run-parts /etc/periodic/5min